# scripts/solix-live.service
[Unit]
Description=Anker Solix MI80 → InfluxDB Live Service
After=network-online.target influxd.service
Wants=network-online.target

[Service]
Type=simple
# Working directory contains anker-solix-api for "from api.api import ..."
WorkingDirectory=/home/nils/shelly/anker-solix-api
# Project directory with uv pyproject.toml
ExecStart=/home/nils/.local/bin/uv run \
    --project /home/nils/shelly/anker-solix-influx \
    solix-live-service

# Environment configuration
EnvironmentFile=/home/nils/shelly/anker-solix-influx/.env

Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=solix-live

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
